# src/jmf.service
[Unit]
Description=Run jmf-service
RefuseManualStart=no

[Service]
Type=simple
ExecStart=/usr/bin/jmf-service
