# src/jmf.timer
[Unit]
Description=jmf.service every 15 minutes

[Timer]
OnBootSec=100
OnUnitActiveSec=15min
Unit=jmf.service

[Install]
WantedBy=timers.target
