# src/what.service
[Unit]
Description=Run what-service
RefuseManualStart=no

[Service]
Type=simple
ExecStart=/usr/bin/what-service
